# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create pository user and group
        if ! getent group pository > /dev/null 2>&1; then
            addgroup --system pository
        fi
        
        if ! getent passwd pository > /dev/null 2>&1; then
            adduser --system --ingroup pository --home /var/lib/pository \
                    --no-create-home --shell /usr/sbin/nologin pository
        fi
        
        # Create directories
        mkdir -p /var/lib/pository
        mkdir -p /var/log/pository
        mkdir -p /etc/pository
        
        # Set ownership
        chown -R pository:pository /var/lib/pository
        chown -R pository:pository /var/log/pository
        chown -R root:pository /etc/pository
        chmod 750 /etc/pository
        
        # Initialize API keys file if not exists
        if [ ! -f /etc/pository/api-keys.json ]; then
            echo '{"keys":[]}' > /etc/pository/api-keys.json
            chown pository:pository /etc/pository/api-keys.json
            chmod 640 /etc/pository/api-keys.json
        fi
        
        # Install npm dependencies
        cd /usr/share/pository
        npm install --production --omit=dev
        
        # Generate initial admin key if not configured
        if ! grep -q "adminKey:" /etc/pository/config.yaml 2>/dev/null && \
           [ -z "${POSITORY_ADMIN_KEY:-}" ]; then
            ADMIN_KEY=$(head -c 32 /dev/urandom | base64 | tr -d '/+=' | head -c 32)
            echo ""
            echo "=========================================="
            echo "POSITORY ADMIN KEY (save this securely):"
            echo "$ADMIN_KEY"
            echo "=========================================="
            echo ""
            echo "# Generated during installation" >> /etc/pository/config.yaml
            echo "adminKey: $ADMIN_KEY" >> /etc/pository/config.yaml
        fi
        
        # Enable and start service
        systemctl daemon-reload
        systemctl enable pository.service
        systemctl start pository.service || true
        ;;
esac

exit 0
